# systemd/user/foo.timer
[Unit]
Description=A Test of the system timers

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
