# systemd/user/foo.service
[Unit]
Description=A Test of the timers

[Service]
ExecStart=
ExecStart=/home/nigel/.config/systemd/user/atest
