# etc/systemd/system/pacmanUpdate.service
[Unit]
Description=Updates Pacmans Repos

[Service]
ExecStart=
ExecStart=pacman -Sy
